# ping-monitor.service
[Unit]
Description=Ping Monitor Daemon with Telegram Notifications
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/bots/ping-monitor
# Запуск строго через интерпретатор виртуального окружения
ExecStart=/root/bot/ping-monitor/venv/bin/python ping_monitor.py
# Если скрипт упадет по непредвиденной ошибке — systemd его поднимет через 10 сек
Restart=always
RestartSec=10
# Перенаправляем логи в journald
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
